# service-files/citysearch-socket.socket
[Socket]
ListenStream=/run/citysearch-socket.sk

[Install]
WantedBy=sockets.target
